# deploy/aiqa-webapp.optional.service
# Note: This service file is optional. The recommended approach is to use nginx's standard service.
# 
# RECOMMENDED: Use standard nginx service with site config:
# sudo cp deploy/aiqa-webapp.nginx.conf /etc/nginx/sites-available/webapp
# sudo ln -s /etc/nginx/sites-available/webapp /etc/nginx/sites-enabled/
# sudo nginx -t
# sudo systemctl reload nginx
#
# ALTERNATIVE: Use this custom service (requires creating your own nginx.conf):
# 1. Create /opt/aiqa/webapp/nginx.conf with full nginx config
# 2. sudo cp deploy/aiqa-webapp.optional.service /etc/systemd/system/
# 3. sudo systemctl daemon-reload
# 4. sudo systemctl enable aiqa-webapp
# 5. sudo systemctl start aiqa-webapp

[Unit]
Description=AIQA Webapp (Nginx)
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/aiqa/webapp
ExecStart=/usr/sbin/nginx -g 'daemon off;' -c /opt/aiqa/webapp/nginx.conf
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aiqa-webapp

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

